# crontab
0 4 * * * /usr/local/bin/python3 /app/main.py >> /app/cron.log 2>&1
